# docker/etc/supervisor/conf.d/dpt.conf
[program:dpt-backend]
command = /usr/src/app/dpt_app/run.sh
directory = /usr/src/app/dpt_app

[program:nginx]
command = /usr/sbin/nginx -g 'daemon off;'
